# tweak/layout/Library/LaunchDaemons/com.barfie.receiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>com.barfie.receiver</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/libexec/receiver</string>
		<string>5000</string>
		<string>"/bin/sh /usr/libexec/receiver.sh 1"</string>
		<string>"/bin/sh /usr/libexec/receiver.sh 0"</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
</dict>
</plist>
